# roles/storage/templates/homelab-disk-space-alert.service.j2
[Unit]
Description=Homelab disk-space alert check

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/homelab-disk-space-alert
